# src/tracker-extract/tracker-extract.service.in
[Unit]
Description=Tracker metadata extractor
ConditionUser=!root

[Service]
Type=dbus
BusName=org.freedesktop.Tracker3.Miner.Extract
ExecStart=@libexecdir@/tracker-extract-3
Restart=on-failure
# Don't restart after tracker daemon -k (aka tracker-control -k)
# Don't restart after seccomp raises SIGSYS either
RestartPreventExitStatus=SIGKILL SIGSYS
Slice=background.slice

[Install]
WantedBy=default.target
