# src/control/tracker-control.service.in
[Unit]
Description=Tracker Filesystem Indexing Proxy

[Service]
Type=dbus
BusName=org.freedesktop.Tracker3.Miner.Files.Control
ExecStart=@libexecdir@/localsearch-control-3
Restart=on-failure
# Don't restart after tracker daemon -k (aka tracker-control -k)
RestartPreventExitStatus=SIGKILL
