# src/indexer/tracker-miner-fs.service.in
[Unit]
Description=LocalSearch indexer
ConditionUser=!@system
ConditionEnvironment=XDG_SESSION_CLASS=user
After=gnome-session.target

[Service]
Type=notify
BusName=org.freedesktop.LocalSearch3
ExecStart=@libexecdir@/localsearch-3
Restart=on-failure
# Don't restart after tracker daemon -k (aka tracker-control -k)
RestartPreventExitStatus=SIGKILL
Slice=background.slice
